# minus.service
[Unit]
Description=Minus - HDMI Ad Blocker
Documentation=https://github.com/garagehq/minus
After=multi-user.target
Conflicts=gdm.service gdm3.service lightdm.service sddm.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/radxa/Minus

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=OPENCV_LOG_LEVEL=ERROR
Environment=PYTORCH_MATCHER_LOGLEVEL=WARNING
Environment=TORCH_LOGS=-all
Environment=TRANSFORMERS_VERBOSITY=error
# Use FastVLM venv site-packages for VLM dependencies (ml_dtypes, transformers, axengine)
Environment=PYTHONPATH=/home/radxa/axera_models/FastVLM-0.5B/fastvlm05_venv/lib/python3.11/site-packages:/home/radxa/.local/lib/python3.11/site-packages

# Pre-start: ensure clean state
# Stop display managers (gdm3, etc.) to release DRM/KMS resources
ExecStartPre=/bin/bash -c 'systemctl stop gdm3 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'systemctl stop gdm 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'systemctl stop lightdm 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'systemctl stop sddm 2>/dev/null || true'
# Kill any lingering Xorg processes
ExecStartPre=/bin/bash -c 'pkill -9 Xorg 2>/dev/null || true'
# Clean up stale processes
ExecStartPre=/bin/bash -c 'pkill -9 ustreamer 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'pkill -9 -f "python3.*[m]inus.py" 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'fuser -k /dev/video0 2>/dev/null || true'
# Clean up stale frame files
ExecStartPre=/bin/bash -c 'rm -f /dev/shm/minus_frame*.jpg /dev/shm/minus_vlm_frame*.jpg 2>/dev/null || true'
# Switch to VT1 and wait for cleanup
ExecStartPre=/bin/bash -c 'chvt 1 2>/dev/null || true'
# 5 second delay before starting to help with stabilization
ExecStartPre=/bin/sleep 5

# Main process
ExecStart=/usr/bin/python3 /home/radxa/Minus/minus.py

# Graceful stop
ExecStop=/bin/kill -TERM $MAINPID
ExecStopPost=/bin/bash -c 'pkill -9 ustreamer 2>/dev/null || true'
ExecStopPost=/bin/bash -c 'rm -f /dev/shm/minus_frame*.jpg /dev/shm/minus_vlm_frame*.jpg 2>/dev/null || true'

# Restart policy
Restart=on-failure
RestartSec=5

# Logging - use journal (view with: journalctl -u minus -f)
StandardOutput=journal
StandardError=journal

# Resource limits
Nice=-5
IOSchedulingClass=realtime
IOSchedulingPriority=0

[Install]
WantedBy=multi-user.target
