# stream-sentry.service
[Unit]
Description=Stream Sentry - HDMI Ad Blocker
Documentation=https://github.com/garagehq/stream-sentry
After=multi-user.target
Conflicts=gdm.service gdm3.service lightdm.service sddm.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/radxa/StreamSentry

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=OPENCV_LOG_LEVEL=ERROR

# Pre-start: ensure clean state
# Stop display managers (gdm3, etc.) to release DRM/KMS resources
ExecStartPre=/bin/bash -c 'systemctl stop gdm3 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'systemctl stop gdm 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'systemctl stop lightdm 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'systemctl stop sddm 2>/dev/null || true'
# Kill any lingering Xorg processes
ExecStartPre=/bin/bash -c 'pkill -9 Xorg 2>/dev/null || true'
# Clean up stale processes
ExecStartPre=/bin/bash -c 'pkill -9 ustreamer 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'pkill -9 -f stream_sentry 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'fuser -k /dev/video0 2>/dev/null || true'
# Clean up stale frame files
ExecStartPre=/bin/bash -c 'rm -f /dev/shm/stream_sentry_frame*.jpg /dev/shm/stream_sentry_vlm_frame*.jpg 2>/dev/null || true'
# Switch to VT1 and wait for cleanup
ExecStartPre=/bin/bash -c 'chvt 1 2>/dev/null || true'
ExecStartPre=/bin/sleep 2

# Main process
ExecStart=/usr/bin/python3 /home/radxa/StreamSentry/stream_sentry.py

# Graceful stop
ExecStop=/bin/kill -TERM $MAINPID
ExecStopPost=/bin/bash -c 'pkill -9 ustreamer 2>/dev/null || true'
ExecStopPost=/bin/bash -c 'rm -f /dev/shm/stream_sentry_frame*.jpg /dev/shm/stream_sentry_vlm_frame*.jpg 2>/dev/null || true'

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=5

# Logging
StandardOutput=append:/tmp/stream_sentry.log
StandardError=append:/tmp/stream_sentry.log

# Resource limits
Nice=-5
IOSchedulingClass=realtime
IOSchedulingPriority=0

[Install]
WantedBy=multi-user.target
